# automation/energy-review-publisher.service
[Unit]
Description=Publish ready Energy Handbook Codex Cloud reviews
After=network-online.target energy-review-executor.service
Wants=network-online.target

[Service]
Type=simple
User=energy-review
Group=energy-review
WorkingDirectory=/opt/energy-review-executor
EnvironmentFile=/etc/energy-review-executor.env
ExecStart=/opt/energy-review-executor/bin/node /opt/energy-review-executor/automation/review-publisher.mjs
Restart=on-failure
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/energy-review
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
CapabilityBoundingSet=
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
